# crontab
# Run every Monday at 3 AM
00 3 * * Mon root python /PlexEmail/scripts/plexEmail.py --config /config/config.conf;

